# packaging/systemd/mouthwrite.service
[Unit]
Description=MouthWrite Linux Voice Input Daemon
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mouthwrite-linux start
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info,mouthwrite_linux=debug

[Install]
WantedBy=default.target
